# panorama-agent.service
# Exemplo de unit systemd para rodar o agente como serviço no Linux.
# Instalação:
#   1. sudo cp panorama-agent.service /etc/systemd/system/
#   2. Ajustar WorkingDirectory e o token abaixo
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now panorama-agent

[Unit]
Description=Panorama Home Center - Agente de Monitoramento
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/panorama-agente
Environment=PANORAMA_API_URL=https://painel.panoramahomecenter.com.br/api
Environment=PANORAMA_AGENT_TOKEN=SUBSTITUA_PELO_TOKEN_DO_ATIVO
Environment=PANORAMA_INTERVAL_MS=60000
ExecStart=/usr/bin/node /opt/panorama-agente/agente-exemplo.js
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
